# templates/com.openai.proxy.antigravity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openai.proxy.antigravity</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>__INSTALL_DIR__/dist/server.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>3000</string>
        <key>LOG_LEVEL</key>
        <string>info</string>
        <key>ADMIN_ENABLED</key>
        <string>true</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/openai-proxy-antigravity.log</string>
    
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/openai-proxy-antigravity-error.log</string>
    
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
